# symbols-database-protector.service
[Unit]
Description=My Symbols Database Protection Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=dansidanutz
Group=staff
WorkingDirectory=/Users/dansidanutz/Desktop/ZmartBot/backend/zmart-api
ExecStart=/usr/bin/python3 /Users/dansidanutz/Desktop/ZmartBot/backend/zmart-api/protect_symbols_database.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Protection settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/Users/dansidanutz/Desktop/ZmartBot/backend/zmart-api

[Install]
WantedBy=multi-user.target
